# ansible/roles/robocoop/templates/robocoop.service.j2
[Unit]
Description=RoboCoop

[Service]
Type=simple
ExecStart=/usr/bin/python2 {{ install_dir }}/robocoop.py --direction auto
WorkingDirectory={{ install_dir }}
